# cmdpadd-client/debian/postinst
#!/bin/sh
#Ironically, we're using sysv here.
invoke-rc.d dbus reload

#DEBHELPER#
